# omnibus/files/server-ctl-cookbooks/infra-server/templates/default/opc_logrotate.cron
#!/bin/sh

PATH=/opt/<%= ChefUtils::Dist::Org::LEGACY_CONF_DIR %>/embedded/sbin:$PATH

command -v logrotate >/dev/null 2>&1 || exit 0
nice -n 19 ionice -c3 logrotate -s /var/log/<%= ChefUtils::Dist::Org::LEGACY_CONF_DIR %>/logrotate.status /etc/<%= ChefUtils::Dist::Org::LEGACY_CONF_DIR %>/logrotate.conf
