# dev/cookbooks/provisioning/files/default/dvm-fs-sync.conf
description "A quick and dirty service to keep guest files synced with host"

start on runlevel [2345]
stop on runlevel [!2345]
kill signal SIGTERM
respawn

exec "/etc/init.d/dvm-fs-sync start"


