# files/private-chef-cookbooks/private-chef/templates/default/opc_log_cleanup.cron
# rendered by private-chef-ctl reconfigure, do not modify
MAILTO=<%= node['private_chef']['notification_email'] %>
53 3 * * * root <%= @cron_command %>
