# files/private-chef-cookbooks/private-chef/templates/default/opc_logrotate.cron
MAILTO=<%= node['private_chef']['notification_email'] %>
53 4 * * * root /usr/sbin/logrotate /etc/opscode/logrotate.conf
