# files/private-chef-cookbooks/private-chef/templates/default/opc_logrotate.cron
#!/bin/sh

test -x logrotate || exit 0
nice -n 19 ionice -c3 logrotate /etc/opscode/logrotate.conf
