# files/private-chef-cookbooks/private-chef/templates/default/opc_logrotate.cron
#!/bin/sh

/usr/bin/test -x /usr/sbin/logrotate || exit 0
/usr/bin/nice -n 19 /usr/bin/ionice -c3 /usr/sbin/logrotate /etc/opscode/logrotate.conf
